# packaging/systemd/opsctl-install-check.timer
[Unit]
Description=Run opsctl installation layout check daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
